# misc/crontab
* * * * * DISPLAY=:0 /home/simon/bin/status
