# crontab
* * * * * /bin/who | /bin/grep --quiet $USER && DISPLAY=:0 /home/$USER/bin/status
